# image/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/php/sbin/php-fpm

[program:snmpd]
command=/usr/sbin/snmpd -LS0-6d -f

[program:spine]
command=/usr/local/spine/bin/spine
